# orientationpro.service
[Unit]
Description=Orientation Pro Congo Application
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/orientationpro
ExecStart=/opt/orientationpro/start-app.sh start
ExecStop=/opt/orientationpro/start-app.sh stop
ExecReload=/opt/orientationpro/start-app.sh restart
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target 